# wps/usr/conf/supervisor/init.d/php-fpm.ini
[program:php-fpm]
user = wordpress
command = /usr/bin/php-fpm -c /home/wordpress/conf/php/php.ini -y /home/wordpress/conf/php/php-fpm.conf
autostart = true
autorestart = true
stdout_events_enabled = true
stderr_events_enabled = true
priority = 200

[program:php-fpm-log]
command=tail -f /home/wordpress/logs/php/php-fpm.log
stdout_events_enabled=true
stderr_events_enabled=true
priority = 1200
