# wps/usr/config/supervisor/init.d/php-fpm.ini
[program:php-fpm]

command = /usr/bin/php-fpm -c /home/wordpress/config/php/php.ini -y /home/wordpress/config/php/php-fpm.conf

autostart = true
autorestart = true
user = wordpress
