# wps/etc/init.d/nginx.ini
[program:nginx]
command = sudo /usr/sbin/nginx -c /home/wordpress/conf.d/nginx.conf
autostart = false
autorestart = true
redirect_stderr = true
stdout_logfile = /home/wordpress/log/access.log
stderr_logfile = /home/wordpress/log/error.log
stopsignal = QUIT