# wps/etc/init.d/nginx.ini
[program:nginx]

command = /usr/sbin/nginx -c /home/wordpress/etc/nginx.conf

autostart = true
autorestart = true
