# wps/etc/init.d/nginx.ini
[program:nginx]
command = /usr/sbin/nginx -c /etc/nginx/nginx.conf -g "error_log /home/wordpress/log/nginx/error.log;"
autostart = true
autorestart = true
