# wps/etc/init.d/nginx.ini
[program:nginx]
command = sudo /usr/sbin/nginx -c /home/wordpress/conf.d/wordpress.conf
autostart = true
autorestart = true
